# Linux/Files/hopperpw.conf
[program:hopperpw]
command = nohup watch -n 60 curl -s http://urlaqui.hopper.pw:ADaSddAS@www.hopper.pw/nic/update
                    ; Command to start app
user = pi                                                          ; User to run as
stdout_logfile = ~/hopperpw.log   ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                       ; Set UTF-8 as default encoding
